# .docker/supervisor/master.ini

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
logfile_maxbytes=50MB        
logfile_backups=10           
loglevel=info                
pidfile=/tmp/supervisord.pid 
minfds=1024                 
minprocs=200                

[program:crond]
command=crond -f -L 15
user=root
autostart=true


[program:api]
directory=/var/www/app
command=yarn start
startsecs=1
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[include]
files = /etc/supervisor.d/*.ini

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file
